# conf/supervisord_blog.conf
[program:blog]
user = root
command=/usr/local/bin/uwsgi --ini /www/blog/conf/uwsgi.ini
autostart=true
autorestart=true
stderr_logfile = /www/log/blog_err.log
stdout_logfile = /www/log/blog_out.log
